# google calendar/cron/crontab
#crontab
29 18 * * * root *applications/AK/modules/sc.py
